# etc/cron.d/coenzyme.umcg-gd-dm.cron
0 7 * * *           /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load ConcordanceCheck/1.0.0; copyConcordanceCheckData.sh -g umcg-gd"
