# etc/cron.d/leucine-zipper.umcg-gd-ateambot.cron
*/5 * * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load ConcordanceCheck/1.0.0; notifications.sh -g umcg-gd -e"
0 1 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load ConcordanceCheck/1.0.0; ConcordanceMakeSamplesheet.sh -g umcg-gd -a umcg-gap"
0 2 * * * /bin/bash -c "export SOURCE_HPC_ENV="True"; . ~/.bashrc; module load ConcordanceCheck/1.0.0; ConcordanceCheck.sh -g umcg-gd"
